# scripts/systemd/slm-polkit-agent.service
[Unit]
Description=SLM Polkit Authentication Agent (Bootstrap)
PartOf=graphical-session.target
After=graphical-session.target

[Service]
ExecStart=%h/Development/Qt/Slm_Desktop/build/slm-polkit-agent
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
