# scripts/systemd/system/slm-post-update-verify.service
[Unit]
Description=SLM Post-Update Verification
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/lib/slm-package-policy/post-update-verify.sh

[Install]
WantedBy=multi-user.target
