# scripts/systemd/system/slm-post-update-verify.path
[Unit]
Description=Trigger SLM Post-Update Verification on dpkg status changes

[Path]
PathChanged=/var/lib/dpkg/status
Unit=slm-post-update-verify.service

[Install]
WantedBy=multi-user.target
